# src/server/cron/jo-mm-m1.cron
#
#   jo-mm-m1.cron
#


#
#   file sync
#
59 23 * * 0 sh /Volumes/JO002TSE001/JO002TSE001/jo-drive/dev/git/c-labs/src/server/sh/jo-freefilesync.sh >> /Volumes/JO002TSE001/JO002TSE001/jo-drive/dev/git/c-labs/src/server/log/jo-freefilesync-`date +\%Y\%m\%d`-`date +\%H\%M\%S`.log 2>&1

#
#   test
#
#* * * * * sh /Volumes/JO002TSE001/JO002TSE001/jo-drive/dev/git/c-labs/src/server/sh/jo-test.sh >> /Volumes/JO002TSE001/JO002TSE001/jo-drive/dev/git/c-labs/src/server/log/jo-test-`date +\%Y\%m\%d`-`date +\%H\%M\%S`.log 2>&1